# admin/sudo/files/sudo.init
#!/bin/sh /etc/rc.common
# Copyright (C) 2014 OpenWrt.org

START=99

start() {
	[ -d /var/lib/sudo ] || mkdir -m 0700 -p /var/lib/sudo
}
